# uwb_localize.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/run_uwb_localize.sh

[Install]
WantedBy=default.target